# install/supervisord.d/tracker.ini
[program:filebt-tracker]
command=/usr/bin/python /opt/filebt/murder_tracker.py --dfile /opt/filebt/data/tracker-data-file --port 8998
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/filebt/logs/filebt-tracker.log
