# install/supervisord.d/tracker.ini
[program:tracker]
command=/opt/opsbt/murder_tracker.py --dfile /opt/opsbt/data --port 8998
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsbt/logs/tracker.log
